# metamon/ops/upstart/supervisor.conf
description "Start services under supervisor"

start on vagrant-ready or runlevel [2345]
stop on runlevel [!2345]

respawn

console log

script
  sudo supervisorctl start all
end script
