# wordpress/ops/upstart/wp_consul_template.conf
description "Consul"

start on vagrant-ready or runlevel [2345]
stop on runlevel [!2345]

respawn

console log

script
  exec /usr/bin/consul-template \
   -consul 127.0.0.1:8500 \
   -template "/ops/templates/wp-config.ctmpl:/var/www/html/wordpress/wp-config.php: service apache2 restart" >> /var/log/ctemplate_wp.log 2>&1
end script
