# PythonAPI-MongoDB/ops/upstart/consul_client.conf
description "Consul"

start on vagrant-ready or runlevel [2345]
stop on runlevel [!2345]

respawn

console log

script
  if [ -f "/etc/service/consul" ]; then
    . /etc/service/consul
  fi

  # Make sure to use all our CPUs, because Consul can block a scheduler thread
  export GOMAXPROCS=`nproc`

  exec /usr/bin/consul agent \
   -data-dir="/tmp/consul" \
   -config-dir="/etc/consul.d" \
   -retry-join=172.31.62.30 \
   -retry-join=172.31.60.43 \
   -retry-join=172.31.51.157 \
   ${CONSUL_FLAGS} \
   >>/var/log/consul.log 2>&1
end script